# scripts/vllm/openmake-vllm-exaone.service
[Unit]
Description=OpenMake vLLM EXAONE 4.5 Instance
After=network.target

[Service]
Type=simple
User=openmake
EnvironmentFile=/etc/openmake/litellm.env
WorkingDirectory=/opt/openmake_llm
ExecStart=/opt/openmake_llm/scripts/vllm/exaone-serve.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
